# deploy/native/sky-sri-worker.service
[Unit]
Description=Sky SRI Worker (native host)
After=network-online.target sky-sri-api.service
Wants=network-online.target

[Service]
Type=simple
User=b2b
Group=b2b
WorkingDirectory=/home/b2b/apps/sky-sriScrapping
EnvironmentFile=/home/b2b/apps/sky-sriScrapping/.env
Environment=PYTHONPATH=/home/b2b/apps/sky-sriScrapping
ExecStart=/bin/bash -lc 'source /home/b2b/apps/sky-sriScrapping/.venv-native/bin/activate && if [ "${HEADLESS:-true}" = "true" ]; then exec python -m api.app.worker; else exec xvfb-run -a python -m api.app.worker; fi'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
